# tools/auto-zhipin/com.redbook.boss-daily-apply.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.redbook.boss-daily-apply</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/proerror/Documents/redbook/tools/auto-zhipin/scripts/daily_apply.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>3</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>14</integer>
            <key>Minute</key>
            <integer>3</integer>
        </dict>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/proerror/.nvm/versions/node/v24.11.1/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/proerror</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/proerror/Documents/redbook/tools/auto-zhipin/logs/launchd-apply.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/proerror/Documents/redbook/tools/auto-zhipin/logs/launchd-apply-err.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
